# tasks/seeds/csharp-default-precedence/files/deploy/moonrelay.service
[Service]
Environment=Relay__DrainTimeoutSeconds=20
# Incident-response flags must be able to supersede the normal service settings.
ExecStart=/opt/moonrelay/MoonRelay --Relay:MaxConcurrency=8
