# scripts/seed-sync.service
[Unit]
Description=Sync new sol-code seeds into moonshiner, commit and push
# Canonical copy: moonshiner repo, scripts/seed-sync.service.
# Install: cp scripts/seed-sync.{service,timer} ~/.config/systemd/user/

[Service]
Type=oneshot
ExecStart=/var/home/venom/moonshiner/scripts/sync_seeds.sh
