# scripts/seed-sync.service
[Unit]
Description=Sync a configured seed source into Moonshiner
# Canonical copy: moonshiner repo, scripts/seed-sync.service.
# Install: cp scripts/seed-sync.{service,timer} ~/.config/systemd/user/

[Service]
Type=oneshot
EnvironmentFile=%h/.config/moonshiner/seed-sync.env
ExecStartPre=%h/moonshiner/scripts/update_seed_source.sh
ExecStart=%h/moonshiner/scripts/sync_seeds.sh
