# scripts/seed-sync.timer
[Unit]
Description=Run the sol-code seed sync every 5 minutes
# Canonical copy: moonshiner repo, scripts/seed-sync.timer.
# Survives reboots: user lingering is enabled (loginctl enable-linger),
# so the user manager starts this at boot without a login.

[Timer]
OnBootSec=2min
OnStartupSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
