# scripts/seed-sync.service
[Unit]
Description=Sync new sol-code seeds into moonshiner, commit and push
# Canonical copy: moonshiner repo, scripts/seed-sync.service.
# Install: cp scripts/seed-sync.{service,timer} ~/.config/systemd/user/

[Service]
Type=oneshot
Environment=MOONSHINER_SEED_SOURCE=%h/sol-code/tasks/seeds
ExecStart=%h/moonshiner/scripts/sync_seeds.sh
