# tasks/seeds/bash-it-systemdrepair/files/units/catalog-db.service
[Unit]
Description=Catalog database readiness
After=network-online.target

[Service]
Type=oneshot
ExecStart=bin/dependency-ready catalog-db.service
RemainAfterExit=yes
