# tasks/seeds/bash-it-systemdrepair/files/units/telemetry-sidecar.service
[Unit]
Description=Unrelated telemetry sidecar
After=network-online.target

[Service]
Type=oneshot
ExecStart=bin/dependency-ready telemetry-sidecar.service
RemainAfterExit=yes
